# src/bin/lqos_scheduler.service.example
[Unit]
After=lqosd.service

[Service]
WorkingDirectory=/opt/libreqos/src
ExecStart=/usr/bin/python3 /opt/libreqos/src/scheduler.py
Restart=always

[Install]
WantedBy=default.target
