# src/bin/lqos_node_manager.service.example
[Unit]
After=network.service lqosd.service
Requires=lqosd.service

[Service]
WorkingDirectory=/opt/libreqos/src/bin
ExecStart=/opt/libreqos/src/bin/lqos_node_manager
Restart=always
#Turn on debuging for service
#Envirnoment=RUST_LOG=info

[Install]
WantedBy=default.target
