# src/bin/lqos_setup.service.example
[Unit]
Wants=network-online.target
After=network-online.target systemd-resolved.service

[Service]
WorkingDirectory=/opt/libreqos/src/bin
ExecStart=/opt/libreqos/src/bin/lqos_setup web
Restart=on-failure

[Install]
WantedBy=multi-user.target
