# src/bin/lqosd.service.example
[Unit]
After=network.service

[Service]
WorkingDirectory=/opt/libreqos/src/bin
ExecStart=/opt/libreqos/src/bin/lqosd
Restart=always

[Install]
WantedBy=default.target
