# src/bin/lqos_node_manager.service.example
[Unit]
After=network.service

[Service]
WorkingDirectory=/opt/libreqos/src/bin
ExecStart=/opt/libreqos/src/bin/lqos_node_manager
Restart=always

[Install]
WantedBy=default.target
